# etc/pcscd.service.in
[Unit]
Description=PC/SC Smart Card Daemon
Requires=pcscd.socket

[Service]
ExecStart=@sbindir_exp@/pcscd --foreground
ExecReload=@sbindir_exp@/pcscd --hotplug
StandardOutput=syslog

[Install]
WantedBy=multi-user.target
Also=pcscd.socket
