# etc/pcscd.service.in
[Unit]
Description=PC/SC Smart Card Daemon
Requires=pcscd.socket

[Service]
ExecStart=@sbindir_exp@/pcscd --foreground --auto-exit
ExecReload=@sbindir_exp@/pcscd --hotplug

[Install]
Also=pcscd.socket
