# etc/pcscd.service.in
[Unit]
Description=PC/SC Smart Card Daemon
Requires=pcscd.socket
Documentation=man:pcscd(8)

[Service]
ExecStart=@sbindir_exp@/pcscd --foreground --auto-exit $PCSCD_ARGS
ExecReload=@sbindir_exp@/pcscd --hotplug
Environment=PCSCD_ARGS=
EnvironmentFile=-@sysconfdir@/default/pcscd

[Install]
Also=pcscd.socket
